# deploy/cuizhao-scrape.cron
# Run the cuizhao Threads scraper every 15 minutes as ubuntu.
# flock prevents overlap if a run takes longer than the interval.
# Installed by deploy/deploy-monitor.sh to /etc/cron.d/cuizhao-scrape.
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=""
*/15 * * * * ubuntu flock -n /tmp/cuizhao-scrape.lock -c 'cd /home/ubuntu/cuizhao && docker compose -f deploy/docker-compose.yml --env-file .env run --rm scraper >> /home/ubuntu/cuizhao/scrape.log 2>&1'
